# voice-assistant.service
[Unit]
Description=Voice Assistant Daemon
After=sound.target

[Service]
ExecStart=%h/voice-assistant/.venv/bin/python %h/voice-assistant/assistant_daemon.py
WorkingDirectory=%h/voice-assistant
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
