# voxlyn-ai.service
[Unit]
Description=voxlyn-ai: Hands-free AI voice assistant daemon
After=sound.target

[Service]
ExecStart=%h/Workspace/voxlyn-ai/.venv/bin/python %h/Workspace/voxlyn-ai/daemon.py
WorkingDirectory=%h/Workspace/voxlyn-ai
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
